# print-server.service
[Unit]
Description=Print Server Dashboard
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=cupsadmin
WorkingDirectory=/home/cupsadmin/print-track
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 